# systemd/system/colameta-mcp-remote.service
[Unit]
Description=ColaMeta public OAuth MCP origin
Documentation=file:/home/jenn/src/colameta-dev/docs/private-beta-systemd.md
After=network-online.target colameta-stable.service
Wants=network-online.target colameta-stable.service
PartOf=colameta-private-beta.target
StartLimitIntervalSec=60s
StartLimitBurst=5
ConditionFileIsExecutable=/home/jenn/tools/colameta/.venv/bin/colameta

[Service]
Type=simple
User=jenn
Group=jenn
WorkingDirectory=/home/jenn/tools/colameta
Environment=HOME=/home/jenn
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=PYTHONUNBUFFERED=1
Environment=MCP_EXPOSURE_PROFILE=commander
Environment=PATH=/home/jenn/.local/bin:/home/jenn/.nvm/versions/node/v22.23.1/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/bin
ExecStart=/home/jenn/tools/colameta/.venv/bin/colameta serve /home/jenn/src/colameta-dev --no-web --mcp-host 127.0.0.1 --mcp-port 8767 --auth-mode external-oauth --public-base-url https://colameta-mcp.skmt617.top --oauth-issuer https://dev-2n3z8xing6eekyok.us.auth0.com/ --oauth-jwks-url https://dev-2n3z8xing6eekyok.us.auth0.com/.well-known/jwks.json --oauth-audience https://colameta-mcp.skmt617.top/mcp --oauth-scopes mcp:read,mcp:preview --no-register-selected
ExecStartPost=/usr/bin/curl --fail --silent --show-error --output /dev/null --retry 30 --retry-all-errors --retry-delay 1 --max-time 5 http://127.0.0.1:8767/healthz
Restart=always
RestartSec=3s
TimeoutStartSec=60s
TimeoutStopSec=30s
KillSignal=SIGINT
KillMode=mixed
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
SystemCallArchitectures=native
TasksMax=512
LimitNOFILE=8192
LogNamespace=colameta

[Install]
WantedBy=multi-user.target
